# templates/service/ipset.timer.j2
[Unit]
Description=Refresh ipsets

[Timer]
{# TODO: randomize? #}
OnCalendar={{ iptables_ipset_refresh_timer }}
Persistent=true

[Install]
WantedBy=timers.target
